# aarch64-buildroot-linux-gnu/sysroot/etc/profile.d/mali-priority.sh
export MALI_SCHED_RT_THREAD_PRIORITY=95
